# openclaude.service
# OpenClaude systemd user service
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp openclaude.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable openclaude.service
#   systemctl --user start openclaude.service
#
# Note: Adjust ExecStart to match the Python that has claude-agent-sdk installed.
#   Check with: python3 -c "import claude_agent_sdk; print('ok')"
#   If using pyenv: ExecStart=/home/is0383kk/.pyenv/shims/python3 -m src.daemon

[Unit]
Description=OpenClaude Resident AI Agent Daemon
After=default.target

[Service]
Type=simple
ExecStart=python3 -m src.daemon
WorkingDirectory=/home/is0383kk/.openclaude
Restart=on-failure
RestartSec=5
KillMode=process
StandardOutput=append:/home/is0383kk/.openclaude/daemon.log
StandardError=append:/home/is0383kk/.openclaude/daemon.log

[Install]
WantedBy=default.target
